# vendor/github.com/phocs/vde_plug_docker/vde_plug_docker.service
[Unit]
Description=VDE network plugin for Docker
Documentation=https://github.com/vixpi/vde_plug_docker
After=docker.service
Requires=docker.service

[Service]
ExecStart=/usr/local/bin/vde_plug_docker
Restart=on-failure

[Install]
WantedBy=multi-user.target
